# deploy/systemd-user/openclaw-gateway.service.d/20-auto-fix-openclaw.conf
[Unit]
OnFailure=auto-fix-openclaw-onfailure.service
StartLimitIntervalSec=60
StartLimitBurst=5
