# deploy/systemd-user/auto-fix-openclaw.timer
[Unit]
Description=Run auto-fix-openclaw every minute

[Timer]
OnBootSec=1min
OnUnitActiveSec=1min
RandomizedDelaySec=8s
Persistent=true
Unit=auto-fix-openclaw.service

[Install]
WantedBy=timers.target
